# config/os-hardening/dslv-zpdi.service
[Unit]
Description=DSLV-ZPDI Production Pipeline (SPEC-011)
After=network.target chrony.service dslv-zpdi-preflight.service
Wants=chrony.service
Requires=dslv-zpdi-preflight.service

[Service]
Type=simple
User=dynogator
Group=dynogator
WorkingDirectory=/home/dynogator/dslv-zpdi
ExecStart=/home/dynogator/dslv-zpdi/venv/bin/python -m dslv_zpdi.main_pipeline
Environment=PYTHONUNBUFFERED=1
Environment=DSLV_CONFIG_PATH=/home/dynogator/dslv-zpdi/config/deployment.yaml
Restart=on-failure
RestartSec=30
StartLimitInterval=120
StartLimitBurst=3
ExecStartPre=/bin/sleep 3
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=4
CPUAffinity=2 3
MemoryMax=2G
MemorySwapMax=0
LimitMEMLOCK=512M
LimitNOFILE=65536

# Systemd Sandboxing
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
NoNewPrivileges=true
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=false
RestrictNamespaces=false

# Read-write paths needed for pipeline operation
ReadWritePaths=/home/dynogator/dslv-zpdi/output /var/lib/dslv_zpdi /run/dslv-zpdi /dev/bus/usb

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dslv-zpdi
RuntimeDirectory=dslv-zpdi
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
