# config/dslv-zpdi-node-receiver.service
[Unit]
Description=DSLV-ZPDI Mobile Node Telemetry Receiver (SPEC-014)
# Accepts JSON telemetry POSTs from Pixel 9 Pro XL / GrapheneOS and other
# swarm nodes over the PiRepo hotspot (10.42.0.x).  Also hosts the RadonEye
# Pro SPEC-015 staging endpoint.
After=network.target dslv-zpdi-hackrf-init.service dslv-zpdi-preflight.service
Wants=dslv-zpdi-preflight.service

[Service]
Type=simple
User=dynogator
Group=dynogator
WorkingDirectory=/home/dynogator/Gem-home/dslv-zpdi
ExecStart=/home/dynogator/Gem-home/dslv-zpdi/venv/bin/python -m dslv_zpdi.layer3_telemetry.node_receiver
Environment=PYTHONUNBUFFERED=1
Environment=DSLV_PRIMARY_OUTPUT_DIR=/home/dynogator/Gem-home/dslv-zpdi/output/primary
Environment=DSLV_SECONDARY_OUTPUT_DIR=/home/dynogator/Gem-home/dslv-zpdi/output/secondary
Environment=DSLV_RECEIVER_PORT=5775
Restart=on-failure
RestartSec=15

ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
NoNewPrivileges=true
ReadWritePaths=/home/dynogator/Gem-home/dslv-zpdi/output

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dslv-zpdi-node-receiver

[Install]
WantedBy=multi-user.target
