# config/dslv-zpdi-baseline.service
[Unit]
Description=DSLV-ZPDI Baseline Capture (SPEC-009.1)
After=network.target chrony.service

[Service]
Type=simple
User=dynogator
Group=dynogator
WorkingDirectory=/home/dynogator/Desktop/DSLV-ZPDI_GitHub_Dev/dslv-zpdi
ExecStart=/home/dynogator/Desktop/DSLV-ZPDI_GitHub_Dev/dslv-zpdi/.venv/bin/python /home/dynogator/Desktop/DSLV-ZPDI_GitHub_Dev/dslv-zpdi/tools/capture_baseline.py
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=20
ExecStartPre=/bin/sleep 3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dslv-zpdi-baseline

[Install]
WantedBy=multi-user.target
